# jellyfin
#!/sbin/openrc-run

: ${JELLYFIN_CACHE_DIR:=/var/cache/jellyfin}
: ${JELLYFIN_LOG_DIR:=/var/log/jellyfin}
: ${JELLYFIN_DATA_DIR:=/var/lib/jellyfin}
: ${JELLYFIN_CONFIG_DIR:=/etc/jellyfin}
: ${JELLYFIN_WEB_DIR:=/usr/share/jellyfin/web}
: ${JELLYFIN_FFMPEG_DIR:=/usr/lib/jellyfin-ffmpeg}

pidfile="/run/${RC_SVCNAME}.pid"
command="/usr/bin/jellyfin"
command_args="--logdir ${JELLYFIN_LOG_DIR} --cachedir ${JELLYFIN_CACHE_DIR} --configdir ${JELLYFIN_CONFIG_DIR} --datadir ${JELLYFIN_DATA_DIR} --webdir ${JELLYFIN_WEB_DIR} --ffmpeg ${JELLYFIN_FFMPEG_DIR}/ffmpeg"
command_user="${RC_SVCNAME}"
command_background=true

start_pre() {
	# Ensure that our dirs are correct
	checkpath --directory --owner jellyfin:jellyfin --mode 0774 \
		${JELLYFIN_LOG_DIR}
	checkpath --directory --owner jellyfin:jellyfin --mode 0774 \
		${JELLYFIN_CACHE_DIR}
	checkpath --directory --owner jellyfin:jellyfin --mode 0774 \
		${JELLYFIN_CONFIG_DIR}
	checkpath --directory --owner jellyfin:jellyfin --mode 0774 \
		${JELLYFIN_DATA_DIR}
	checkpath --directory --owner jellyfin:jellyfin --mode 0774 \
                ${JELLYFIN_WEB_DIR}
	checkpath --directory --owner jellyfin:jellyfin --mode 0774 \
                ${JELLYFIN_FFMPEG_DIR}
}
